# debian/radeon-profile-daemon.postrm
#!/bin/sh
set -e

echo "Disabling/removing radeon profile daemon symlink"
systemctl disable radeon-profile-daemon.service
wait $!

exit 0
